# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/mongodb/bin/mongod --port 27017 --dbpath /data/db
autostart=true
autorestart=true
startsecs=2
user=root

